# epoll-server.service
[Unit]
Description=epoll-server
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/epoll-server
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
